# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/config.nginx
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:era_test_node]
command=/usr/local/bin/era_test_node run
autostart=true
autorestart=true
stdout_logfile=/var/log/era_test_node.log
stderr_logfile=/var/log/era_test_node.log
